# config/local.gpt-qq-bot.chat-hub.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.gpt-qq-bot.chat-hub</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__PROJECT_DIR__" || exit 1; [ -f "config/local.env" ] &amp;&amp; source "config/local.env"; npm start</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>
  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/runtime/logs/chat-hub.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/runtime/logs/chat-hub.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
